# deploy/systemd/bot-trade.timer
# Systemd timer file for bot-trade
# Install to: /etc/systemd/system/bot-trade.timer

[Unit]
Description=Bot Trade - Trading Schedule Timer
Requires=bot-trade.service

[Timer]
# PUT Session: 12:00, 12:30, 13:00, 13:30 UTC
OnCalendar=*-*-* 12:00:00
OnCalendar=*-*-* 12:30:00
OnCalendar=*-*-* 13:00:00
OnCalendar=*-*-* 13:30:00

# CALL Session: 18:00, 18:30 UTC
OnCalendar=*-*-* 18:00:00
OnCalendar=*-*-* 18:30:00

# Run on boot if missed
Persistent=true

[Install]
WantedBy=timers.target
