# deploy/systemd/bot-trade.timer
# Systemd timer file for bot-trade
# Install to: /etc/systemd/system/bot-trade.timer

[Unit]
Description=Bot Trade - Trading Schedule Timer
Requires=bot-trade.service

[Timer]
# Run every 30 minutes (24/7)
# Bot will check trading_hours and session_filters internally
OnCalendar=*:0/30

# Run on boot if missed
Persistent=true

[Install]
WantedBy=timers.target
